# deploy/systemd/xp-bot.service
[Unit]
Description=XP Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/xp-bot
EnvironmentFile=/opt/xp-bot/.env
ExecStart=/opt/xp-bot/.venv/bin/python /opt/xp-bot/bot.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
